# systemd/dummy.service
[Service]
Type=simple
ExecStart=printf "Hello, world!\n"
